# docs/claude-code/launchd/dev.dopesoft.mcp-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  mcp-proxy bridges Claude Code's stdio MCP transport to SSE on 127.0.0.1:8765.
  cloudflared exposes :8765 publicly behind Cloudflare Access.

  The `--` separates mcp-proxy flags from the inner `claude mcp serve` invocation.
  In MCP-serve mode the parent MCP client (Infinity) is the permission
  authority — Claude Code just exposes tools — so no `--dangerously-skip-permissions`
  is needed (the flag doesn't exist on `claude mcp serve` anyway). Infinity's
  ClaudeCodeGate routes risky calls through mem_trust_contracts.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.dopesoft.mcp-proxy</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>exec mcp-proxy --port 8765 --host 127.0.0.1 --server sse -- claude mcp serve</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- __NPM_BIN__ is the npm global prefix (e.g. ~/.npm-global/bin)
             where mcp-proxy and claude actually live. install.sh resolves it
             via `npm prefix -g`/bin so the LaunchAgent finds both binaries. -->
        <key>PATH</key>
        <string>__NPM_BIN__:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>__HOME__/Dev</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>NetworkState</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/infinity-coder/mcp-proxy.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/infinity-coder/mcp-proxy.err</string>
</dict>
</plist>
